# opentsdb/supervisor.conf
[supervisord]
nodaemon=true

[program:hbase]
command=/opt/hbase-service.sh
stopwaitsecs=60

[program:opentsdb]
command=/opt/opentsdb-service.sh
stopwaitsecs=60
